# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx.stdout.log
stderr_logfile=/var/log/nginx.stderr.log

[program:api]
command=node /app/api/server.js
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/api.stdout.log
stderr_logfile=/var/log/api.stderr.log
environment=NODE_ENV="production"
